# systemd/gpsd.service.d/override.conf
[Service]
EnvironmentFile=/etc/default/gpsd
Restart=always
RestartSec=5
